# deploy/systemd/muhasebe-ops-backup.timer
# /etc/systemd/system/muhasebe-ops-backup.timer
# Muhasebe Operasyonları — Daily backup timer (taslak, Faz 17B)
# Bu fazda enable/start YAPILMAZ. Production deploy fazında:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now muhasebe-ops-backup.timer

[Unit]
Description=Daily backup timer for Muhasebe Operasyonlari
Documentation=file:/var/www/muhasebe-ops/repo/_docs/BACKUP_RESTORE_PLAN.md
Requires=muhasebe-ops-backup.service

[Timer]
# Her gün 03:30 yerel saat
OnCalendar=*-*-* 03:30:00
# Sistem uyuyor/kapalıysa ilk açılışta kaçırılan job çalışsın
Persistent=true
# Çakışma azaltma: 0–10 dk rastgele gecikme
RandomizedDelaySec=600
Unit=muhasebe-ops-backup.service

[Install]
WantedBy=timers.target
